# files/giving-day-socket-server.conf
description "Giving Day node.js server"

start on started mountall
stop on shutdown

# automatically respawn
respawn
respawn limit 99 5

script
    export HOME="/root"
    exec /usr//bin/nodejs /etc/nodejs/giving-day/server.js >> /var/log/node.log 2>&1
end script
